# supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:api]
directory=/app/
command=./main
autostart=true
autorestart=true
stderr_logfile=/app/api_err.log
stdout_logfile=/app/api_out.log